# hadoop/profile.d/hadoop.sh
#!/usr/bin/env bash

export HADOOP_HOME=/hadoop
export YARN_HOME=$HADOOP_HOME
